# autostart/qt5-polkit-agent.desktop
[Desktop Entry]
Type=Application
Exec=qt5-polkit-agent
TryExec=qt5-polkit-agent

Name=PolicyKit Agent
